# tunnel/service.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon for cloudflared on the Mac mini gateway.

  Installed at /Library/LaunchDaemons/ai.owera.cloudflared.plist (chown root:wheel,
  chmod 644), loaded with:

      sudo launchctl bootstrap system /Library/LaunchDaemons/ai.owera.cloudflared.plist
      sudo launchctl enable system/ai.owera.cloudflared

  System-wide (not per-user) so the tunnel survives logout. Restart-on-fail is
  handled via KeepAlive with SuccessfulExit=false. The template variable
  conventions mirror owera-fleet/templates/launchd/*.tmpl:
    {{.Label}}      -> "ai.owera.cloudflared"
    {{.HomeDir}}    -> "/var/log/cloudflared"  (logs root, not $HOME)
    {{.ConfigPath}} -> "/usr/local/etc/cloudflared/config.yml"
    {{.BinPath}}    -> "/opt/homebrew/bin/cloudflared" or "/usr/local/bin/cloudflared"

  This file is the *example* checked into git. The rendered copy lives on the
  gateway and is regenerated by the install script.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{.Label}}</string>

    <key>UserName</key>
    <string>_cloudflared</string>
    <key>GroupName</key>
    <string>_cloudflared</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{.BinPath}}</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>{{.ConfigPath}}</string>
        <string>--no-autoupdate</string>
        <string>run</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TUNNEL_METRICS</key>
        <string>127.0.0.1:9300</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>{{.HomeDir}}/cloudflared.out.log</string>
    <key>StandardErrorPath</key>
    <string>{{.HomeDir}}/cloudflared.err.log</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
